# systemd/cf-ddns.timer
[Unit]
Description=Runs Dynamic DDNS on a Timer

[Timer]
OnBootSec=10min
OnUnitActiveSec=1min

[Install]
WantedBy=timers.target
